# sporty.socket
[Unit]
Description=sporty socket

[Socket]
ListenStream=/run/sporty.sock

[Install]
WantedBy=sockets.target